# newscrux.service
# Newscrux AI News Push Notification Service
# Edit User=, WorkingDirectory=, EnvironmentFile=, and --lang= to match your setup.
# Example Ubuntu setup:
#   User=ubuntu
#   WorkingDirectory=/home/ubuntu/newscrux
#   EnvironmentFile=/home/ubuntu/newscrux/.env
#   ExecStart=/usr/bin/node /home/ubuntu/newscrux/dist/index.js --lang=uz

[Unit]
Description=Newscrux AI News Push Notification Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/newscrux
ExecStart=/usr/bin/node /home/ubuntu/newscrux/dist/index.js --lang=en
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=5

EnvironmentFile=/home/ubuntu/newscrux/.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=newscrux

NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/ubuntu/newscrux/data

[Install]
WantedBy=multi-user.target
